# ai-gf-webhook.service
[Unit]
Description=AI Girlfriend Bot - Webhook Server (YooKassa Payments)
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/ai_gf
Environment="PATH=/root/ai_gf/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/ai_gf/venv/bin/python webhook_server.py

# Автоматический перезапуск при падении
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-gf-webhook

[Install]
WantedBy=multi-user.target
